# systemd/couch-nagger.service
[Unit]
Description=Couch Nagger - Dog Detection Service
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/couch-nagger

# Environment
Environment="PATH=/home/pi/couch-nagger/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Load config from file if exists
EnvironmentFile=-/home/pi/couch-nagger/config/production.env

# Run the monitoring script
ExecStart=/home/pi/couch-nagger/.venv/bin/couch-monitor

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=couch-nagger

# Resource limits (optional - adjust as needed)
# MemoryMax=512M
# CPUQuota=50%

[Install]
WantedBy=multi-user.target
